# pkg/stage/usr/local/etc/rc.d/darwin
#!/bin/sh
# PROVIDE: darwin
# KEYWORD: shutdown
# This file is part of Vulture 3.
#
# Vulture 3 is free software: you can redistribute it and/or modify
# it under the terms of the GNU General Public License as published by
# the Free Software Foundation, either version 3 of the License, or
# (at your option) any later version.

# Vulture 3 is distributed in the hope that it will be useful,
# but WITHOUT ANY WARRANTY; without even the implied warranty of
# MERCHANTABILITY or FITNESS FOR A PARTICULAR PURPOSE.  See the
# GNU General Public License for more details.

# You should have received a copy of the GNU General Public License
# along with Vulture 3.  If not, see http://www.gnu.org/licenses/.

########################## WARNING ###############################
# Don't edit this file, it is automatically generated by Vulture #
########################## WARNING ###############################

. /etc/rc.subr

name="darwin"
rcvar="${name}_enable"

load_rc_config "$name"
: ${darwin_log_level="WARNING"}
: ${darwin_enable="YES"}
: ${darwin_user="darwin"}
: ${darwin_group="vlt-conf"}
: ${darwin_config="/home/darwin/conf/$name.conf"}

procname="/home/vlt-os/env/bin/python"
command="/usr/sbin/daemon"
pidfile="/var/run/darwin/manager.pid"
command_args="-f -p ${pidfile} ${procname} /home/darwin/manager/manager.py ${darwin_config} -l ${darwin_log_level}"
required_files="${darwin_config}"

start_postcmd="def_set_unix_perms"


def_set_unix_perms() {
    /bin/sleep 2
    /usr/sbin/chown $darwin_user:$darwin_group /var/sockets/darwin/*.sock 2> /dev/null
    /bin/chmod -h 775 /var/sockets/darwin/*.sock 2> /dev/null
}

extra_commands="reload"

run_rc_command "$1"
